# etc/cron.d/pcompact
# Copyright (c) 1999-2017, Parallels International GmbH
# Compacting virtual disks from 2am to 4am
0 2 * * * root /usr/sbin/pcompact -t 2h -q >/dev/null
