# etc/cron.d/pcompact
# Copyright (C) Parallels IP Holdings GmbH, 1999-2012.
# Compacting virtual disks from 2am to 4am
0 2 * * * root /usr/sbin/pcompact -t 2h
